# com.user.file-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly Auto-Organizer — launchd config
  
  Setup:
    1. Copy this file to ~/Library/LaunchAgents/
    2. Replace YOUR_USERNAME with your macOS username
    3. Run: launchctl load ~/Library/LaunchAgents/com.user.file-organizer.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.file-organizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/YOUR_USERNAME/file-organizer/auto-organize.sh</string>
    </array>

    <!-- Run every Sunday at 10:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.file-organizer/logs/launchd-error.log</string>
</dict>
</plist>
